# target/linux/ath79/tiny/base-files/etc/hotplug.d/ieee80211/10_fix_wifi_mac
#!/bin/ash

[ "$ACTION" = "add" ] || exit 0

PHYNBR=${DEVPATH##*/phy}

[ -n $PHYNBR ] || exit 0

. /lib/functions.sh
. /lib/functions/system.sh
. /lib/functions/uci-defaults.sh
. /lib/functions/k2t.sh

board=$(board_name)

case "$board" in
	dlink,dap-1120-a1|\
	dlink,dap-1320-ax|\
	dlink,dap-1320-bx|\
	dlink,dap-1320-cx)
		mtd_get_mac_text "art" 0xffa0 > /sys${DEVPATH}/macaddress
		ucidef_set_label_macaddr $(mtd_get_mac_text "art" 0xffa0)
		;;
esac
